# deploy/garrett-api.service
[Unit]
Description=Garrett Intelligence Hub FastAPI
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/crewai-agent
ExecStart=/root/crewai-agent/venv/bin/uvicorn api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=OTEL_SDK_DISABLED=true

[Install]
WantedBy=multi-user.target
